# markdoc.cron
*/5  *   *   *   *   cd /site && markdoc build
